# root/etc/s6-overlay/s6-rc.d/init-setupcron/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

SCRIPT_NAME="INIT-SETUPCRON"

log "$SCRIPT_NAME" "Starting cron configuration"

# Validate cron-related environment variables
if [ -n "${recreate_vpn_cron:-}" ]; then
    log "$SCRIPT_NAME" "VPN recreation cron: $recreate_vpn_cron ($(parse_cron "$recreate_vpn_cron"))"
else
    log "$SCRIPT_NAME" "VPN recreation cron not configured"
fi

if [ -n "${check_connection_cron:-}" ]; then
    log "$SCRIPT_NAME" "Health check cron: $check_connection_cron ($(parse_cron "$check_connection_cron"))"
    
    # Validate CHECK_CONNECTION_* variables that are used by vpn-healthcheck
    
    if [ -z "$check_connection_attempts" ] || [ "$check_connection_attempts" -le 0 ]; then
        log_error "$SCRIPT_NAME" "ERROR: CHECK_CONNECTION_ATTEMPTS must be a positive integer, got: '$check_connection_attempts'"
        exit 1
    fi
    
    if [ -z "$check_connection_url" ]; then
        log_error "$SCRIPT_NAME" "ERROR: CHECK_CONNECTION_URL must be set"
        exit 1
    fi
    
    if [ -z "$check_connection_attempt_interval" ] || [ "$check_connection_attempt_interval" -le 0 ]; then
        log_error "$SCRIPT_NAME" "ERROR: CHECK_CONNECTION_ATTEMPT_INTERVAL must be a positive integer, got: '$check_connection_attempt_interval'"
        exit 1
    fi
    
    # Validation passed
else
    log "$SCRIPT_NAME" "Health check cron not configured, skipping validation"
fi

cron_dir="/var/spool/cron/crontabs"
cron_file="$cron_dir/root"

rm -f "$cron_file"
touch "$cron_file"

if [ -n "${recreate_vpn_cron:-}" ]; then
    echo "$recreate_vpn_cron vpn-reconnect" >> "$cron_file"
else
    log "$SCRIPT_NAME" "VPN reconnection cron not configured"
fi

if [ -n "${check_connection_cron:-}" ]; then
    echo "$check_connection_cron vpn-healthcheck" >> "$cron_file"
else
    log "$SCRIPT_NAME" "Health check cron not configured"
fi

log "$SCRIPT_NAME" "Cron configuration completed"
exit 0
