# root/etc/s6-overlay/s6-rc.d/init-createvpnconfig/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

vpn-config

exit 0
