# root/etc/s6-overlay/s6-rc.d/svc-nordvpn/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

SCRIPT_NAME="SERVICE-NORDVPN"

log "$SCRIPT_NAME" "Starting NordVPN service"

# Select a server and write the config (blocks forever on fatal error so s6
# keeps the container alive), then bring up the tunnel via the shared helper.
vpn-config
apply_wireguard_config "$SCRIPT_NAME" || { log_error "$SCRIPT_NAME" "CRITICAL: could not apply WireGuard config"; sleep infinity; }

# Wait for VPN connection
log "$SCRIPT_NAME" "Waiting for VPN connection..."
counter=0
while ! is_vpn_connected && [ $counter -lt 30 ]; do
    sleep 2
    counter=$((counter + 1))
done

if is_vpn_connected; then
    log "$SCRIPT_NAME" "VPN connected successfully"
    # Run network diagnostic if enabled
    if [ "$network_diagnostic_enabled" = "true" ]; then
        log "$SCRIPT_NAME" "Running network diagnostic..."
        network-diagnostic
    fi
else
    log "$SCRIPT_NAME" "VPN connection timeout"
fi

# Keep the s6 service running indefinitely
sleep infinity
