# root/etc/s6-overlay/s6-rc.d/nordvpnd/finish
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

echo "[SERVICE-NORDVPND] Stopping nordvpnd service..."

wgfile="/etc/wireguard/wg0.conf"

. /usr/local/bin/backend-functions

if wg show wg0 >/dev/null 2>&1; then
    wg-quick down wg0 || true
else
    echo "[SERVICE-NORDVPND] Interface wg0 not found, skipping wg-quick down"
fi
rm -f "$wgfile"

echo "[SERVICE-NORDVPND] Cleaning VPN pinhole..."
run4 -F VPN-SERVER

echo "[SERVICE-NORDVPND] nordvpnd service stopped successfully"
exit 0
