# root/etc/s6-overlay/s6-rc.d/crond/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

echo "[SERVICE-CROND] Starting crond service..."

exec crond -f
