# root/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

SCRIPT_NAME="SERVICE-CRON"

log "$SCRIPT_NAME" "Starting cron service"

exec crond -f
